# etc/supervisord.d/hbase_regionserver.ini
[program:hbase_regionserver]
command=/opt/hbase/bin/start_regionserver.sh
stopwaitsecs=120
prioirty=10
autorestart=true
startsecs=4
