# etc/supervisord.d/hbase_master.ini
[program:hbase_master]
command=/opt/hbase/bin/start_master.sh
stopwaitsecs=120
prioirty=5
autorestart=true
