# etc/supervisord.d/zookeeper.ini
[program:zookeeper]
command=/opt/hbase/bin/start_zookeeper.sh
stopwaitsecs=120
prioirty=3
autorestart=true
